# release/titus-isolate.service
[Unit]
Description=Titus container isolation
After=docker.service
Requires=docker.service

[Service]
ExecStart=/usr/share/python/titus-isolate/bin/python3 /usr/share/python/titus-isolate/bin/main.py --package-count 2 --cores-per-package 16
Restart=always
RestartSec=5

[Install]
Alias=titus-isolate
